# packages/sys-auth/nss-pam-ldapd/files/systemd/nslcd.service
[Unit]
Description=nslcd

[Service]
User=nslcd
Group=nslcd
Type=simple
ExecStart=/usr/host/bin/nslcd --nofork

[Install]
WantedBy=multi-user.target

